# imgproxy.service
[Unit]
Description=Image proxy for appliance-budget-picker
After=network.target

[Service]
User=nick
WorkingDirectory=/home/nick/nickpalm/renopicker
ExecStart=/home/nick/.local/bin/uvicorn imgproxy:app --host 127.0.0.1 --port 8768
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
